# services/watchtower/ansible/roles/daily_summary/templates/daily-summary.service.j2
[Unit]
Description=Little Wolf Acres — Daily Summary
After=network-online.target prometheus.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/daily-summary.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-summary
